# debian/uhub.postrm
#!/bin/sh -e

if [ "$1" = purge ]; then
        update-rc.d uhub remove >/dev/null
fi

