# debian/uhub.prerm
#!/bin/sh -e

if [ "$1" = remove ]; then
        if command -v invoke-rc.d >/dev/null 2>&1; then
                invoke-rc.d uhub stop || true
        else
                /etc/init.d/uhub stop
        fi
fi

